# supervisor_config/hat.conf
[program:hat]
command=/opt/caproto-venv/bin/python /opt/iocs/hat.py --pins 11,16 -v --prefix=hat{{ prefix }}:
