# supervisor_config/rgb.conf
[program:rgb]
command=/opt/caproto-venv/bin/python /opt/iocs/rgb.py --red 15 --green 22 --blue 11 -vvv --list-pvs
